# init/powerman.service
# ~/.config/systemd/user/powerman.service
[Unit]
Description=Powerman

[Service]
ExecStart=/usr/bin/powerman

[Install]
WantedBy=default.target

